# home/private_dot_config/exact_shellrc/profile.d/02_xdg_browser.tmpl
{{- if .is_wsl }}
export BROWSER=$HOME/.local/bin/wsl-open
{{- end }}
