# home/private_dot_config/exact_shellrc/profile.d/02_xdg_browser.tmpl
{{- if and .is_wsl (lookPath "wslview") }}
export BROWSER=wslview
{{- end }}
